# services/testy-crash@.service
[Unit]
Description="testy-crash: A crash service for the testy framework"
Documentation=https://github.com/wiredtiger/testy

[Service]
User=testy
Group=testy
Type=oneshot
ExecStartPre=/bin/bash -c '! systemctl is-active --quiet testy-backup@%I.service || (echo "Skipping scheduled crash because of a scheduled backup currently in progress." && exit 1)'
ExecStartPre=/bin/bash -c 'if [ $(systemctl show --property MainPID --value testy-run@%I.service) == 0 ]; then echo "No active service matching \'testy-run@%I.service\' found, skipping scheduled crash." && exit 1; fi'
ExecStartPre=/bin/bash -c 'sudo pkill -P $(systemctl show --property MainPID --value testy-run@%I.service) --signal SIGKILL'
ExecStart=/bin/bash -c '${script_dir}/testy-snapshot.sh ${workload_dir}/%I/%I.sh testy-crash'
ExecStartPost=/bin/bash -c 'sudo systemctl start testy-run@%I.service'
TimeoutSec=3600s
StandardOutput=journal+console
StandardError=journal+console
