# services/testy-backup@.service
[Unit]
Description="testy-backup: An AWS snapshot-based database backup and validation service for the testy framework"
Documentation=https://github.com/wiredtiger/testy

[Service]
User=testy
Group=testy
ExecStartPre=/bin/bash -c '(systemctl status testy-crash@%I.service | grep -E "Active: inactive|failed") || (echo "Skipping scheduled backup because of a scheduled crash currently in progress." && exit 1)'
ExecStart=/bin/bash -c '${script_dir}/testy-snapshot.sh ${workload_dir}/%I/%I.sh testy-backup'
TimeoutSec=3600s
StandardOutput=journal+console
StandardError=journal+console
