# services/testy-backup@.service
[Unit]
Description="testy-backup: An AWS snapshot-based database backup and validation service for the testy framework"
Documentation=https://github.com/wiredtiger/testy

[Service]
User=testy
Group=testy
# The crash service is of type oneshot, thus we cannot use the is-active command.
ExecStartPre=/bin/bash -c 'systemctl status testy-crash@%I.service | grep -E "Active: inactive|failed"'
ExecStart=/bin/bash -c '${script_dir}/testy-snapshot.sh ${workload_dir}/%I/%I.sh testy-backup'
TimeoutSec=3600s
StandardOutput=journal+console
StandardError=journal+console
